# debian/templates/etc/systemd/system/genisys-powermate-to-rotator.service
[Unit]
Description=Smartspaces Pi Dashboard
Requires=network.target

[Service]
Type=simple
#POSTINSTALL_USER_HOOK
PermissionsStartOnly=true
ExecStart=/usr/local/bin/genisys-powermate-to-rotator

[Install]
WantedBy=multi-user.target
